# cps-3.5/deb/squeeze/DEBIAN/postinst
#!/bin/sh
ZUSER=zope:zope
CPS_HOME=/opt/cps-3.5
INSTANCE=/var/opt/cps-3.5/test-instance
RC_SCRIPT=zope-cps-test

case "$1" in
  abort-deconfigure|abort-upgrade)
    ;;
  abort-remove)
  ;;
  configure)
    # Byte compiling
    python2.4 $CPS_HOME/bin/compilezpy.py $CPS_HOME/Products > /dev/null
    if [ -d $INSTANCE ]; then
	echo "Found existing instance at $INSTANCE"
	upscript=$INSTANCE/bin/deb_upgrade
	echo "Upgrading via $upscript"
	exec $upscript $2
    else
        echo Creating Zope/CPS instance at $INSTANCE
        $CPS_HOME/bin/mkcpsinstance --posix-owner=$ZUSER -u admin:admin $INSTANCE

	ln -s $CPS_HOME/Products/CPSDefault/themes $INSTANCE/cps-themes
        echo Setting up startup script
        ln -s $INSTANCE/bin/zopectl /etc/init.d/zope-cps-test
        update-rc.d -f zope-cps-test start 92 2 3 4 5 . stop 92 0 1 6 .
    fi
    echo "Starting Zope/CPS test instance (listens on port 6500)"
    /etc/init.d/zope-cps-test start

    ;;

  *)
    echo "postinst called with unknown argument \`$1'" >&2
    exit 1
    ;;
esac

